# root/etc/services.d/fleet/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/bin/java -Dfleet.config.base=/config -jar /app/fleet/fleet.jar
